# templates/consul.ini
[program:consul]
command=consul agent -config-dir=/etc/consul
stopsignal=INT
stdout_logfile=/var/log/consul.log
autostart=true
autorestart=true
